# backup/allspark-backup.service

[Unit]
Description=Sauvegarde incrémentielle

[Service]
Type=simple
ExecStart=/srv/allspark-backup.sh
StandardError=journal
WorkingDirectory=/srv
