# backup/db-backup.timer
[Unit]
Description=Sauvegarde quotidienne des bases de données

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=false

[Install]
WantedBy=timers.target
