# backup/increments-backup.service

[Unit]
Description=Sauvegarde incrémentielle

[Service]
Type=simple
ExecStart=/srv/increments-backup.sh
StandardError=journal
